# neo4j/supervisord.conf

[supervisord]
nodaemon=true

[program:sshd]
command=service ssh start

[program:neo4j]
command=service neo4j-service start
