# Services/neo4j/supervisord.conf

[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:neo4j]
command=service neo4j-service start
